# deploy/systemd/market-sentinel-health.service
[Unit]
Description=MarketSentinel loopback health check
Wants=network-online.target
After=network-online.target market-sentinel-web.service
StartLimitIntervalSec=5min
StartLimitBurst=5

[Service]
Type=oneshot
User=market-sentinel
Group=market-sentinel
WorkingDirectory=/opt/market-sentinel
EnvironmentFile=-/etc/market-sentinel/market-sentinel.env
ExecStart=/opt/market-sentinel/.venv/bin/python /opt/market-sentinel/scripts/verify_service_health.py --retries 2 --retry-delay 5
TimeoutStartSec=30
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectHostname=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
